# Packages/KernBypass Daemon/silica_data/scripts/postinst
#!/bin/sh
chown root:wheel /Library/LaunchDaemons/jp.akusio.kernbypassd.plist
launchctl load /Library/LaunchDaemons/jp.akusio.kernbypassd.plist
